# systemd.service
[Unit]
Description=lume.io
After=network.target

[Service]
User=infamous
Group=infamous
WorkingDirectory=/home/infamous/lume+lume/apps/website/
ExecStart=/usr/local/bin/npm start
Restart=always
StandardOutput=syslog
SyslogIdentifier=lume-website
Environment=BLAH=blah

[Install]
WantedBy=multi-user.target
